# pi/services/birdpi-bootstrap.service
[Unit]
Description=BirdPi one-shot bootstrap
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/birdpi-bootstrap.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
